# rsd/rsd-arm-linux/ROOTFS/var/lib/dpkg/info/resolvconf.conffiles
/etc/network/if-up.d/000resolvconf
/etc/network/if-down.d/resolvconf
/etc/dhcp/dhclient-enter-hooks.d/resolvconf
/etc/resolvconf/interface-order
/etc/resolvconf/resolv.conf.d/base
/etc/resolvconf/resolv.conf.d/head
/etc/resolvconf/update.d/dnscache
/etc/resolvconf/update.d/libc
/etc/init.d/resolvconf
/etc/ppp/ip-up.d/000resolvconf
/etc/ppp/ip-down.d/000resolvconf
